# .configs/etc/rc.local


# wait for all container boot up
sleep 90
# restart nginx after all other container boot finish
cd /opt/demp-stack/nginx && make restart

exit 0
